# app-emulation/ganeti-webmgr/files/ganeti-webmgr-vncproxy
#!/sbin/runscript
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-proxy/wwwoffle/files/wwwoffled,v 1.1 2005/04/22 19:15:12 mrness Exp $

depend() {
	need net
}

start() {
	ebegin "Starting gwm-vncproxy"
	twistd --pidfile=/var/run/vncap.pid vncap -c tcp:8888
	eend $?
}

stop() {
	ebegin "Stopping gwm-vncproxy"
	start-stop-daemon --stop --quiet --pidfile /var/run/vncap.pid 
	eend $?
}
