# app-emulation/ganeti-webmgr/files/ganeti-webmgr
#!/sbin/runscript

extra_commands="migrate"

description_migrate="Run South migration for the database schemas."
description_stop="Kills all children and stops the server."

PIDFILE="/var/run/ganeti_webmgr/gwm.pid"
GWM_PATH="/var/lib/ganeti_webmgr"

# default values
GWM_SERVER_OPTIONS="${GWM_SERVER_OPTIONS:--v 0}"
GWM_SERVER_PORT=${GWM_SERVER_PORT:-}

depend() {
	need net
}

checkconfd() {
	CONFIGFILE="${GWM_PATH}/settings.py"
	if [ ! -r "${CONFIGFILE}" ]; then
		eerror "Unable to find configuration file: ${CONFIGFILE}"
		eerror "You need to create and customize Ganeti Web Manager settings:"
		eerror "# cp ${GWM_PATH}/settings.py.dist ${GWM_PATH}/settings.py"
		return 1
	fi
}

checkconfig() {
	# validating configuration file
	checkconfd || return 1

	# validating settings
	return 0
}

start() {
	checkconfig || return 1

	ebegin "Starting Ganeti web Manager"
	cd "${GWM_PATH}"
	start-stop-daemon --start --quiet -i -b -m \
		--pidfile "${PIDFILE}" \
		--chdir "${GWM_PATH}" \
		--exec /usr/bin/python2.7 \
		  -- manage.py runserver ${GWM_SERVER_OPTIONS} ${GWM_SERVER_PORT}
	eend $? "Failed to start Ganeti web Manager"
}

stop() {
	ebegin "Stopping Ganeti web Manager"
	# stop ganeti-webmgr via pkill as start-stop-daemon does not stop child process.
	pkill -P $( cat "${PIDFILE}" )
	rm -f "${PIDFILE}"
	eend $? "Failed to stop Ganeti web Manager"
}

migrate() {
	# validating configuration file
	checkconfd || return 1

	ebegin "Upgrading database"
	cd "${GWM_PATH}"
	/usr/bin/python2.7 "${GWM_PATH}"/manage.py migrate
	eend $? "Failed to upgrade database"
}
