# app-emulation/ganeti-webmgr/files/ganeti-webmgr-flashpolicy
#!/sbin/runscript
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-proxy/wwwoffle/files/wwwoffled,v 1.1 2005/04/22 19:15:12 mrness Exp $

source /etc/conf.d/gwm

depend() {
	need net
}

start() {
	ebegin "Starting gwm-flashpolicy"

	export PYTHONPATH="$GWM_PATH"
	source $GWM_PATH/bin/activate

	twistd --pidfile=/var/run/flashpolicy.pid flashpolicy
	eend $?
}

stop() {
	ebegin "Stopping gwm-flashpolicy"
	start-stop-daemon --stop --quiet --pidfile /var/run/flashpolicy.pid 
	eend $?
}
